# main_ws/sys/freshclam/int_freshclam.timer
[Unit]
Description=Run freshclam weekly.
Requires=int_freshclam.service

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
