# main_ws/sys/clamscan/int_clamscan.timer
[Unit]
Description=Run clamscan weekly.
Requires=int_clamscan.service

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
